# znc/files/init-znc.conf.jinja
# znc

description "IRC Bouncer"

start on runlevel [2345]

stop on runlevel [016]

respawn
respawn limit 10 5
setuid znc

script
  exec /usr/bin/znc --datadir=/home/znc/.znc -f
end script
